# telegram-2fa-bot.service
[Unit]
Description=Telegram 2FA Bot
After=network.target

[Service]
User=your_user
WorkingDirectory=/home/your_user/telegram-2fa-bot
ExecStart=/usr/bin/python3 /home/your_user/telegram-2fa-bot/bot.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
